# systemd/system/mnt-devname.automount
# Enable this unit `systemctl enable mnt-devname.automount`.
[Unit]
Description=Automount /etc/systemd/system/mnt-devname.mount

[Automount]
Where=/mnt/target

[Install]
WantedBy=multi-user.target

