# etc/systemd/system/netctl@dev.static.service.d/profile.conf
[Unit]
Description=Static ethernet connection for NAME.
BindsTo=sys-subsystem-net-devices-NAME.device
After=sys-subsystem-net-devices-NAME.device
