# etc/systemd/system/mnt-storage.automount
# Enable this unit `systemctl enable mnt-storage.automount`.
[Unit]
Description=Automount /etc/systemd/system/mnt-storage.mount

[Automount]
Where=/mnt/storage

[Install]
WantedBy=multi-user.target
