# etc/systemd/system/srv-nfs.automount
[Unit]
Description=Automount NFS storage
Requires=network-online.target
#After=netctl.static.service # ? not necessary ?

[Automount]
Where=/srv/nfs

[Install]
WantedBy=multi-user.target
