# systemd/system/netctl@dev.static.service.d/profile.conf
[Unit]
Description=Static ethernet connection for eno1.
BindsTo=sys-subsystem-net-devices-eno1.device
After=sys-subsystem-net-devices-eno1.device
